# archiso/airootfs/etc/systemd/system/lnos-autostart.service
[Unit]
Description=LnOS Auto-start Installer
After=multi-user.target
ConditionPathExists=/usr/local/bin/setup-lnos-autostart.sh

[Service]
Type=oneshot
ExecStartPre=/usr/local/bin/setup-lnos-autostart.sh
ExecStart=/usr/local/bin/lnos-autostart.sh
StandardOutput=journal
StandardError=journal
RemainAfterExit=yes
TimeoutSec=300

[Install]
WantedBy=multi-user.target